# xcertcheck@.service
[Unit]
Description=Certificate check service

[Service]
Type=oneshot
ExecStart=/usr/bin/xcertcheck
